# pr2-netconsole/root/etc/cron.d/checknetconsole
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * *  root if [ `cat /sys/kernel/config/netconsole/logging/remote_mac` = 'ff:ff:ff:ff:ff:ff' ]; then service netconsole start; fi
